# Run/supervisor/conf.d/wb_batch_manage.ini
[program:proxy_pool]
command = gunicorn main:run
directory = /usr/share/nginx/proxy_pool/Run/
user = root
autostart = true
autorestart = true
startsecs = 10
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
stdout_logfile = /var/log/proxy_pool/proxy_pool.log
stderr_logfile = /var/log/proxy_pool/proxy_pool_err.log

